# jenkins/docker/etc/service/jenkins/run
#!/bin/bash
export HOME=/srv/jenkins
/bin/tini -- /usr/local/bin/jenkins.sh
