# jenkins/docker/etc/service/jenkins/run
#!/bin/bash
export HOME=/var/jenkins_home
/bin/tini -- /usr/local/bin/jenkins.sh
