# jenkins/docker/etc/service/jenkins/run
#!/bin/bash

exec /sbin/setuser app /usr/local/bin/jenkins.sh
