# jenkins/docker/etc/service/jenkins/run
#!/bin/bash

exec /sbin/setuser app /scripts/base/docker-wrapper-task.sh "opg-jenkins" /usr/local/bin/jenkins.sh
