# ruby-dsd/docker/service/app/run
#!/bin/sh

cd /app
exec su app -c 'ruby hello.rb'
