# nginx0x644/docker/service/nginx/run
#!/bin/sh

mkdir -p /var/log/nginx

chown app:app /var/log/app

exec /usr/sbin/nginx
