# salt-master/docker/service/salt-master/run
#!/bin/sh

/etc/init.d/salt-master start
