# elasticsearch50x644/docker/service/elasticsearch/run
#!/bin/bash
#
export HOME=/usr/share/elasticsearch
chown -R elasticsearch:elasticsearch /usr/share/elasticsearch/data
exec gosu elasticsearch elasticsearch
