# ruby/docker/service/app/run
#!/bin/sh

cd /app
exec su app -c 'ruby hello.rb' &> /var/log/app/app.log

