# devise/docker/service/app/run
#!/bin/sh


su app -c "cd /app && bundle exec rake db:mongoid:create_indexes RAILS_ENV=production"
exec su app -c "cd /app && bundle exec unicorn -E production -o unix:///var/run/app/app.sock"
