# elasticsearch5/docker/service/elasticsearch/run
#!/bin/bash
sysctl -w vm.max_map_count=262144
export HOME=/usr/share/elasticsearch
chown -R elasticsearch:elasticsearch /usr/share/elasticsearch/data
exec gosu elasticsearch elasticsearch
