# monitoring/docker/service/kibana/run
#!/bin/sh

exec /opt/logstash/bin/logstash web -a 0.0.0.0 -p 80
