# kibana/docker/service/kibana/run
#!/bin/bash
exec gosu kibana /usr/share/kibana/bin/kibana
