# kibana/docker/service/kibana/run
#!/bin/bash
export HOME=/opt/kibana
chown -R kibana:kibana $HOME
exec gosu kibana kibana
