# base0x644/docker/service/dnsmasq/run
#!/bin/bash

# let's call dnsmasq but clean all env variables
exec -c  /usr/sbin/dnsmasq -u root -k
