# fake-sqs/docker/service/fake_sqs/run
#!/bin/bash

/usr/local/bin/fake_sqs -p 4568 --no-daemonize
