# rabbitmq/docker/service/rabbitmq/run
#!/bin/bash
export HOME=/var/lib/rabbitmq
exec rabbitmq-server
