# php-fpm/docker/service/php-fpm/run
#!/bin/sh

/sbin/update-php-fpm-environment.sh

#php-fpm is configured to downgrade permissions to user app
exec /usr/sbin/php5-fpm --fpm-config /etc/php5/fpm/pool.d/app.conf
