# salt-master/docker/service/salt-minion/run
#!/bin/sh

/etc/init.d/salt-minion start
